# cookbooks/selenium/templates/default/org.seleniumhq.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string><%=@name%></string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ProgramArguments</key>
    <array>
      <string><%=@exec%></string>
      <% @args.each do |arg| %>
      <string><%=arg.gsub('"', '')%></string>
      <% end %>
    </array>
    <key>ServiceDescription</key>
    <string><%=@name%> service</string>
    <key>StandardErrorPath</key>
    <string>/var/log/selenium/<%=@name%>.log</string>
    <key>StandardOutPath</key>
    <string>/var/log/selenium/<%=@name%>.log</string>
  </dict>
</plist>
