# cookbooks/selenium/templates/default/sysvinit.erb
#!/bin/bash
#
# <%=@name%>
#
# chkconfig: 345 90 25
# description: Selenium <%=@name%> service

# Source function library.
. /etc/init.d/functions

user=<%=@user%>
exec="<%=@exec%>"
args="<%=@args%>"
lockfile="/var/lock/subsys/<%=@name%>"
pidfile="/var/run/<%=@name%>.pid"
logfile="/var/log/<%=@name%>.log"
prog="<%=@name%>"
xdisplay="<%=@xdisplay%>"
port="<%=@port%>"

RETVAL=0

start() {
  echo -n $"Starting $prog: "

  touch $pidfile
  chown $user $pidfile

  touch $logfile
  chown $user $logfile

  /bin/su - $user -c "DISPLAY=\"$xdisplay\" $exec $args >> $logfile 2>&1 & echo \$! > $pidfile"

  sleep 2

  pgrep -fl $prog
  RETVAL=$?
  [ $RETVAL -eq 0 ] && echo_success || echo_failure

  listen=0
  while [ $listen -eq 0 ];
  do
    listen=`netstat -ltn | awk '{print $4}' | grep "${port}$" | wc -l`
    sleep 1
  done

  echo
  [ $RETVAL -eq 0 ] && touch $lockfile
  return $RETVAL
}

stop() {
  echo -n $"Stopping $prog: "
  killproc -p $pidfile $prog
  RETVAL=$?
  echo
  [ $RETVAL -eq 0 ] && rm -f $lockfile $pidfile
  return $RETVAL
}

restart() {
  stop
  start
}


case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  status)
    status -p ${pidfile} ${prog}
    RETVAL=$?
    ;;
  restart)
    restart
    ;;
  *)
    echo $"Usage: $0 {start|stop|restart}"
    exit 1
esac
